# config/systemd/user/autorandr.service
[Unit]
Description=Auto RandR Daemon

[Service]
ExecStart=/usr/bin/x-on-resize --config '/usr/bin/autorandr -c'

[Install]
WantedBy=graphical.target
