# config/systemd/user/console.target
[Unit]
Description=TTY Console
StopWhenUnneeded=true
AllowIsolate=true

Conflicts=graphical.target
