# config/systemd/user/default.target
[Unit]
Description=Defaults
AllowIsolate=true

Wants=desktop.target
