# config/systemd/user/desktop.target
[Unit]
Description=Desktop Environment
AllowIsolate=true

Conflicts=console.target

Requires=graphical.target
After=graphical.target
