# config/systemd/user/dropbox.service
[Unit]
Description=Dropbox Sync Daemon

[Service]
ExecStart=/usr/bin/dropboxd
Nice=15

[Install]
WantedBy=default.target
