# config/systemd/user/dropbox.service
[Unit]
Description=Dropbox sync daemon
After=network.target
Requires=network.target

[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/dropboxd
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=default.target
