# config/systemd/user/dropbox.service
[Unit]
Description=Dropbox sync daemon

[Service]
ExecStart=/opt/dropbox/dropboxd

[Install]
WantedBy=default.target
