# config/systemd/user/gui@.target
[Unit]
Description=Graphical User Interface (display %i)

Wants=xbase@%i.target

[Install]
WantedBy=default.target
