# config/systemd/user/mopidy.service
[Unit]
Description=Mopidy service

[Service]
ExecStart=/usr/bin/mopidy
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=default.target
