# config/systemd/user/mopidy.service
[Unit]
Description=Mopidy

[Service]
ExecStart=/usr/bin/mopidy
Restart=always

[Install]
WantedBy=default.target
