# config/systemd/user/multiplexer.target
[Unit]
Description=Terminal multiplexer
Documentation=man:tmux(1)

[Install]
Alias=default.target
