# config/systemd/user/multiplexer.target
[Unit]
Description=Terminal Multiplexer

StopWhenUnneeded=true

After=dbus.socket
Wants=dbus.socket

[Install]
WantedBy=default.target
