# config/systemd/user/network.target
[Unit]
Description=Network
Requires=network-wait-online.service

[Install]
WantedBy=default.target
