# config/systemd/user/offlineimap.timer
[Unit]
Description=OfflineIMAP Synchronization Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
