# config/systemd/user/offlineimap.service
[Unit]
Description=OfflineIMAP Service

[Service]
Type=simple
ExecStart=/usr/bin/offlineimap -o

[Install]
WantedBy=default.target
