# config/systemd/user/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
Documentation=man:pulseaudio(1)

Requires=sound.target
After=sound.target

[Service]
Type=forking

ExecStart=/usr/bin/pulseaudio --start
ExecStop=/usr/bin/pulseaudio --kill

[Install]
WantedBy=default.target
