# config/systemd/user/terminals.service
[Unit]
Description=X terminals

[Service]
Type=oneshot
ExecStart=/usr/bin/termite
Restart=on-success

[Install]
WantedBy=xinitrc.target
