# config/systemd/user/tmux.service
[Unit]
Description=Terminal Multiplexer Server
Documentation=man:tmux(1)

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/tmux new-session -d -s %H
ExecStop=/usr/bin/tmux kill-server
KillMode=none

[Install]
WantedBy=default.target
