# config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper

[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/nitrogen --restore

[Install]
WantedBy=xinitrc.target
