# config/systemd/user/wm.slice
[Unit]
Description=Window Manager Slice
