# config/systemd/user/wm.target
[Unit]
Description=Window Manager
StopWhenUnneeded=true

Conflicts=console.target

After=xmonad@0.service

[Install]
WantedBy=desktop.target
