# config/systemd/user/xbase@.target
[Unit]
Description=X Server with essential environment and settings (display %i)

StopWhenUnneeded=true

After=x@%i.service

Requires=x@%i.service
