# config/systemd/user/xcape.service
[Unit]
Description=Combine Ctrl+Escape

[Service]
Type=forking
ExecStart=/usr/bin/xcape

[Install]
WantedBy=graphical.target
