# config/systemd/user/xcompmgr.service
[Unit]
Description=X Composite Manager
Documentation=man:xcompmgr(1)

[Service]
ExecStart=/usr/bin/xcompmgr

[Install]
WantedBy=graphical.target
