# config/systemd/user/xinitrc.target
[Unit]
Description=Xinitrc Stuff
Wants=wm.target

[Install]
Alias=default.target
