# config/systemd/user/xmodmap.service
[Unit]
Description=Xmodmap
After=xorg.target
Requires=xorg.target

[Service]
Type=oneshot
ExecStart=/usr/bin/xmodmap %h/.Xmodmap

[Install]
WantedBy=xinitrc.target
