# config/systemd/user/xorg.target
[Unit]
Description=The basic Xorg server
