# config/systemd/user/xorg@.socket
[Unit]
Description=Xorg Socket (display %i)

[Socket]
ListenStream=/tmp/.X11-unix/X%i
