# config/systemd/user/xserver@.target
[Unit]
Description=X Server on display %i
RefuseManualStart=true

PartOf=xorg@%i.service
