# systemd/.config/systemd/user/blueproximity.service
[Unit]
Description=Proximity Detector

[Service]
Type=simple
ExecStart=/usr/bin/start_proximity.sh

[Install]
WantedBy=graphical.target
