# systemd/.config/systemd/user/break.timer
[Unit]
Description=Break Encouragement

[Timer]
OnBootSec=25min
OnUnitActiveSec=25min

[Install]
WantedBy=timers.target
