# systemd/.config/systemd/user/console.target
[Unit]
Description=TTY Console
StopWhenUnneeded=true
Conflicts=graphical.target
