# systemd/.config/systemd/user/default.target
[Unit]
Description=Defaults
AllowIsolate=true
Wants=console.target
