# systemd/.config/systemd/user/eclim.service
[Unit]
Description=eclim daemon

[Service]
Type=simple
ExecStart=/usr/lib/eclipse/eclimd

[Install]
WantedBy=default.target
