# systemd/.config/systemd/user/graphical.target
[Unit]
Description=Graphical Interface
Wants=input-method.service
Conflicts=console.target
After=input-method.service
AllowIsolate=yes
