# systemd/.config/systemd/user/pulseaudio-dlna.service
[Unit]
Description=PulseAudio Chromecast support

[Service]
ExecStart=/usr/bin/pulseaudio-dlna --port 9999 --cover-mode=application --auto-reconnect

[Install]
WantedBy=default.target
