# email/.config/systemd/user/offlineimap.timer
[Unit]
Description=OfflineIMAP Synchronization Timer

[Timer]
OnUnitInactiveSec=15m
Unit=offlineimap.service

[Install]
WantedBy=timers.target
