# compton/.config/systemd/user/compton.service
[Unit]
Description=Compositor for X11
Documentation=man:compton(1)

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=graphical.target
