# calendar/.config/systemd/user/calendar-alert.timer
[Unit]
Description=Calendar Alert Timer

[Timer]
OnBootSec=2min
OnUnitInactiveSec=2min

[Install]
WantedBy=timers.target
