# calendar/.config/systemd/user/calendar-alert.service
[Unit]
Description=Calendar Alert

[Service]
Type=oneshot
ExecStart=/usr/bin/env notify-calendar-event

[Install]
WantedBy=graphical.target
