# calendar/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=Vdirsyncer Synchronization Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
