# calendar/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=Vdirsyncer Synchronization Timer

[Timer]
OnUnitInactiveSec=15m
Unit=vdirsyncer.service

[Install]
WantedBy=timers.target
